# _archief/com.sprankelend.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sprankelend.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/timrijnders/Documents/Sprankelend/start.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/timrijnders/Documents/Sprankelend</string>

    <!-- Start automatisch bij login / na slaapstand -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Herstart automatisch als hij crasht -->
    <key>KeepAlive</key>
    <true/>

    <!-- Logs voor debugging -->
    <key>StandardOutPath</key>
    <string>/Users/timrijnders/Documents/Sprankelend/.server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/timrijnders/Documents/Sprankelend/.server.log</string>
</dict>
</plist>
